# deploy/systemd/jordan-mentor-dispatch.service
[Unit]
Description=Jordan mentor dispatch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/jordan
Environment=JORDAN_LOG_PATH=/opt/jordan/workspace/logs/jordan.jsonl
ExecStart=/opt/jordan/deploy/systemd/run-mentor-dispatch.sh
StandardOutput=append:/opt/jordan/workspace/logs/jordan.jsonl
StandardError=append:/opt/jordan/workspace/logs/jordan.jsonl

[Install]
WantedBy=multi-user.target
