# deploy/systemd/jordan-mentor-dispatch.service
[Unit]
Description=Jordan mentor dispatch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/jordan
ExecStart=/bin/bash -lc 'exec "${JORDAN_HOME:?set JORDAN_HOME in /etc/default/jordan}/deploy/systemd/run-mentor-dispatch.sh"'

[Install]
WantedBy=multi-user.target
